# autostart/dashboard.desktop
[Desktop Entry]
Type=Application
Name=Dashboard
Exec=python3 /home/a/dashboard.py
